# examples/destiny/node-exporter/templates/node_exporter.service.tmpl
{{/*
  systemd-unit node_exporter. Go text/template strict (ADR-010), Вариант B: operator-input
  читается `.input.<name>` напрямую (без passthrough через params.vars), путь бинаря —
  `.vars.bin_path` (file-var). Перечень полей — destiny.yml input:.

  ExecStart собирается условно (флаг-конвенция node_exporter --help): list → per-item
  --no-collector.<name> / --collector.<name>; gomaxprocs>0 → --runtime.gomaxprocs=N;
  пустые string-параметры (fs/netdev-exclude) флаг НЕ добавляют. Перенос строк через `\`.

  ★ Сервис-аккаунт — стабильный system-user (прод-конвенция §2, stateful): textfile-каталог
  переживает рестарты и читается root-сборщиками → нужен фиксированный uid (DynamicUser тут
  непригоден). Под root основной демон НЕ работает.

  ★ systemd-hardening (прод-конвенция §3, «безопасность на первом месте»): полный блок
  беспривилегированного Go-демона. ReadOnlyPaths=.input.textfile_dir — демон только ЧИТАЕТ
  .prom (пишут oneshot-сборщики). MemoryDenyWriteExecute уместен (Go-бинарь без W+X-страниц).
  CapabilityBoundingSet пуст — публичные метрики привилегий не требуют.
*/}}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ .input.user }}
Group={{ .input.group }}
ExecStart={{ .vars.bin_path }} \
    --collector.textfile.directory={{ .input.textfile_dir }} \
    --web.listen-address={{ .input.listen }} \
    --web.telemetry-path={{ .input.web_telemetry_path }} \
    --log.level={{ .input.log_level }} \
    --log.format={{ .input.log_format }}{{ range .input.disabled_collectors }} \
    --no-collector.{{ . }}{{ end }}{{ range .input.enabled_collectors }} \
    --collector.{{ . }}{{ end }}{{ range $name, $opts := .input.collector_options }}{{ range $k, $v := $opts }} \
    --collector.{{ $name }}.{{ $k }}={{ $v }}{{ end }}{{ end }}{{ if gt (int .input.gomaxprocs) 0 }} \
    --runtime.gomaxprocs={{ .input.gomaxprocs }}{{ end }}{{ if .input.fs_mount_points_exclude }} \
    --collector.filesystem.mount-points-exclude={{ .input.fs_mount_points_exclude }}{{ end }}{{ if .input.netdev_device_exclude }} \
    --collector.netdev.device-exclude={{ .input.netdev_device_exclude }}{{ end }}
SyslogIdentifier=node_exporter
Restart=on-failure
RestartSec=5

# systemd-hardening (прод-конвенция §3).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictNamespaces=yes
ProtectProc=invisible
ReadOnlyPaths={{ .input.textfile_dir }}

[Install]
WantedBy=multi-user.target
