# examples/destiny/node-exporter/templates/node-exporter-ipmitool.timer.tmpl
{{/*
  Таймер IPMI-сборщика. Статичен (vars: {}), .tmpl обязателен для core.file.rendered.
  Та же ConditionPathExists=/dev/ipmi0, что и в .service: защита от ручного start без
  BMC-устройства.
*/}}
[Unit]
Description=Run IPMI metrics collection every 1 minute
ConditionPathExists=/dev/ipmi0

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
