# examples/destiny/node-exporter/templates/node_exporter.service.tmpl
{{/*
  systemd-unit node_exporter. Контекст рендера (docs/templating.md §3.2, Вариант B):
    .vars  — destiny-локалы vars.yml (file-vars), доступны шаблону НАПРЯМУЮ как
             `.vars.<file_var>` без passthrough через params.vars: bin_path —
             целевой путь бинаря (input.bin_dir + '/node_exporter', резолвится в
             vars.yml, не operator-input);
    .input — operator-input прохода (Вариант B, ADR-010 §3.2 amendment): user,
             group, listen, textfile_dir; прод-параметры демона: gomaxprocs
             (int), disabled_collectors / enabled_collectors (list<string>),
             collector_options (map<string,map<string,string>>), log_level,
             log_format, web_telemetry_path, fs_mount_points_exclude,
             netdev_device_exclude (string). Читаются `.input.<name>` напрямую,
             без passthrough через params.vars;
    .self  — стабильные факты хоста (здесь не используются).
  Движок — Go text/template + sprig-allowlist, strict-mode (ADR-010).

  ExecStart собирается условно (флаг-конвенция node_exporter --help): list →
  per-item --no-collector.<name> / --collector.<name>; gomaxprocs>0 →
  --runtime.gomaxprocs=N; пустые string-параметры (fs/netdev-exclude) флаг НЕ
  добавляют — node_exporter применяет встроенный дефолт. systemd допускает
  перенос строк ExecStart через `\` — флаги разнесены для читаемости.

  Сервис-аккаунт — стабильный system-user .vars.user (прод-конвенция §2,
  stateful-ветка): textfile-каталог переживает рестарты и читается root-
  сборщиками, поэтому нужен фиксированный uid-владелец (DynamicUser тут
  непригоден — transient-uid не совпал бы с владельцем .prom между запусками).
  Под root основной демон НЕ работает.

  systemd-hardening (прод-конвенция §3, «безопасность на первом месте»): полный
  блок беспривилегированного Go-демона. ReadOnlyPaths=.input.textfile_dir — демон
  только ЧИТАЕТ .prom из textfile-каталога (пишут туда oneshot-сборщики), поэтому
  каталог для него только-чтение. MemoryDenyWriteExecute уместен (Go-бинарь без
  W+X-страниц). CapabilityBoundingSet пуст — публичные метрики ядра/ФС привилегий
  не требуют.
*/}}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ .input.user }}
Group={{ .input.group }}
ExecStart={{ .vars.bin_path }} \
    --collector.textfile.directory={{ .input.textfile_dir }} \
    --web.listen-address={{ .input.listen }} \
    --web.telemetry-path={{ .input.web_telemetry_path }} \
    --log.level={{ .input.log_level }} \
    --log.format={{ .input.log_format }}{{ range .input.disabled_collectors }} \
    --no-collector.{{ . }}{{ end }}{{ range .input.enabled_collectors }} \
    --collector.{{ . }}{{ end }}{{ range $name, $opts := .input.collector_options }}{{ range $k, $v := $opts }} \
    --collector.{{ $name }}.{{ $k }}={{ $v }}{{ end }}{{ end }}{{ if gt (int .input.gomaxprocs) 0 }} \
    --runtime.gomaxprocs={{ .input.gomaxprocs }}{{ end }}{{ if .input.fs_mount_points_exclude }} \
    --collector.filesystem.mount-points-exclude={{ .input.fs_mount_points_exclude }}{{ end }}{{ if .input.netdev_device_exclude }} \
    --collector.netdev.device-exclude={{ .input.netdev_device_exclude }}{{ end }}
SyslogIdentifier=node_exporter
Restart=on-failure
RestartSec=5

# systemd-hardening (прод-конвенция §3).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictNamespaces=yes
ProtectProc=invisible
ReadOnlyPaths={{ .input.textfile_dir }}

[Install]
WantedBy=multi-user.target
