# examples/destiny/node-exporter/templates/node_exporter.service.tmpl
{{/*
  systemd-unit node_exporter. Контекст рендера (docs/templating.md §3.2):
    .vars  — явно переданные core.file.rendered (bin_dir, user, group, listen,
             textfile_dir);
    .self  — стабильные факты хоста (здесь не используются).
  Движок — Go text/template + sprig-allowlist, strict-mode (ADR-010).

  Сервис-аккаунт — стабильный system-user .vars.user (прод-конвенция §2,
  stateful-ветка): textfile-каталог переживает рестарты и читается root-
  сборщиками, поэтому нужен фиксированный uid-владелец (DynamicUser тут
  непригоден — transient-uid не совпал бы с владельцем .prom между запусками).
  Под root основной демон НЕ работает.

  systemd-hardening (прод-конвенция §3, «безопасность на первом месте»): полный
  блок беспривилегированного Go-демона. ReadOnlyPaths=.vars.textfile_dir — демон
  только ЧИТАЕТ .prom из textfile-каталога (пишут туда oneshot-сборщики), поэтому
  каталог для него только-чтение. MemoryDenyWriteExecute уместен (Go-бинарь без
  W+X-страниц). CapabilityBoundingSet пуст — публичные метрики ядра/ФС привилегий
  не требуют.
*/}}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ .vars.user }}
Group={{ .vars.group }}
ExecStart={{ .vars.bin_dir }}/node_exporter --collector.textfile.directory={{ .vars.textfile_dir }} --web.listen-address={{ .vars.listen }}
SyslogIdentifier=node_exporter
Restart=on-failure
RestartSec=5

# systemd-hardening (прод-конвенция §3).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictNamespaces=yes
ProtectProc=invisible
ReadOnlyPaths={{ .vars.textfile_dir }}

[Install]
WantedBy=multi-user.target
