# examples/destiny/node-exporter/templates/node_exporter.service.tmpl
{{/*
  systemd unit for node_exporter. Go text/template strict (ADR-010), Variant B: operator-input
  is read directly via `.input.<name>` (no passthrough via params.vars), the binary path is
  `.vars.bin_path` (file-var). The field list lives in destiny.yml input:.

  ExecStart is assembled conditionally (the node_exporter --help flag convention): list -> per-item
  --no-collector.<name> / --collector.<name>; gomaxprocs>0 -> --runtime.gomaxprocs=N;
  empty string parameters (fs/netdev-exclude) do NOT add the flag. Line continuation via `\`.

  ★ Service account - a stable system-user (prod convention §2, stateful): the textfile directory
  survives restarts and is read by root-owned collectors -> needs a fixed uid (DynamicUser is
  unsuitable here). The main daemon does NOT run as root.

  ★ systemd-hardening (prod convention §3, "security first"): the full hardening block for an
  unprivileged Go daemon. ReadOnlyPaths=.input.textfile_dir - the daemon only READS
  .prom files (written by oneshot collectors). MemoryDenyWriteExecute is appropriate (a Go binary has no W+X pages).
  CapabilityBoundingSet is empty - publishing metrics requires no privileges.
*/}}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ .input.user }}
Group={{ .input.group }}
ExecStart={{ .vars.bin_path }} \
    --collector.textfile.directory={{ .input.textfile_dir }} \
    --web.listen-address={{ .input.listen }} \
    --web.telemetry-path={{ .input.web_telemetry_path }} \
    --log.level={{ .input.log_level }} \
    --log.format={{ .input.log_format }}{{ range .input.disabled_collectors }} \
    --no-collector.{{ . }}{{ end }}{{ range .input.enabled_collectors }} \
    --collector.{{ . }}{{ end }}{{ range $name, $opts := .input.collector_options }}{{ range $k, $v := $opts }} \
    --collector.{{ $name }}.{{ $k }}={{ $v }}{{ end }}{{ end }}{{ if gt (int .input.gomaxprocs) 0 }} \
    --runtime.gomaxprocs={{ .input.gomaxprocs }}{{ end }}{{ if .input.fs_mount_points_exclude }} \
    --collector.filesystem.mount-points-exclude={{ .input.fs_mount_points_exclude }}{{ end }}{{ if .input.netdev_device_exclude }} \
    --collector.netdev.device-exclude={{ .input.netdev_device_exclude }}{{ end }}
SyslogIdentifier=node_exporter
Restart=on-failure
RestartSec=5

# systemd-hardening (prod convention §3).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictNamespaces=yes
ProtectProc=invisible
ReadOnlyPaths={{ .input.textfile_dir }}

[Install]
WantedBy=multi-user.target
