# examples/destiny/node-exporter/templates/node-exporter-ipmitool.timer.tmpl
{{/*
  Таймер IPMI-сборщика. Контекст рендера пуст (vars: {}) — шаблон статичен,
  .tmpl обязателен для core.file.rendered. Движок — Go text/template strict
  (ADR-010).

  Та же ConditionPathExists=/dev/ipmi0 в [Unit], что и в .service: защита от
  ручного start без BMC-устройства.
*/}}
[Unit]
Description=Run IPMI metrics collection every 1 minute
ConditionPathExists=/dev/ipmi0

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
