# examples/destiny/node-exporter/templates/node-exporter-ipmitool.service.tmpl
{{/*
  oneshot collector of IPMI sensors for node_exporter. Go text/template strict (ADR-010).
  .vars.textfile_dir - directory of the ipmitool awk script AND the output ipmitool.prom.

  Pipeline: `ipmitool sensor` → ipmitool awk script (into prometheus form) → ipmitool.prom.
  ★ User=root: /usr/bin/ipmitool requires access to /dev/ipmi0 (char-ipmidev).

  ★ ProtectKernelModules=no (per spec): the ipmi_si / ipmi_devintf driver may load lazily
  on first access to /dev/ipmi0 - disallowing module loading would break collection on hosts
  where the driver is not yet loaded. The rest of the sandbox is strict.

  ConditionPathExists=/dev/ipmi0 - without a BMC device the collector does not start (same
  Condition in .timer).
*/}}
[Unit]
Description=Collect IPMI metrics for node_exporter
ConditionPathExists=/dev/ipmi0

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/env bash -c 'set -euo pipefail; tmp="$(mktemp {{ .vars.textfile_dir }}/ipmitool.prom.XXXXXX)"; /usr/bin/ipmitool sensor | {{ .vars.textfile_dir }}/ipmitool > "$tmp"; mv "$tmp" "{{ .vars.textfile_dir }}/ipmitool.prom"'

# systemd hardening for a privileged oneshot (per spec, collectors section).
ReadWritePaths={{ .vars.textfile_dir }}
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
ProtectKernelTunables=yes
ProtectKernelModules=no
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX
NoNewPrivileges=yes
ProtectProc=invisible
PrivateDevices=no
DevicePolicy=closed
DeviceAllow=char-ipmidev rw
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_DAC_OVERRIDE
