# examples/destiny/node-exporter/templates/node-exporter-smartmon.timer.tmpl
{{/*
  Таймер SMART-сборщика. Статичен (vars: {}), .tmpl обязателен для core.file.rendered.
  Та же ConditionVirtualization=no, что и в .service: защита от ручного start без железа.
*/}}
[Unit]
Description=Run SMART metrics collection every 1 minute
ConditionVirtualization=no

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
