# examples/destiny/node-exporter/templates/node-exporter-smartmon.timer.tmpl
{{/*
  Таймер SMART-сборщика. Контекст рендера пуст (vars: {}) — шаблон статичен,
  .tmpl обязателен для core.file.rendered. Движок — Go text/template strict
  (ADR-010).

  Та же ConditionVirtualization=no в [Unit], что и в .service: защита от ручного
  start таймера там, где железа нет.
*/}}
[Unit]
Description=Run SMART metrics collection every 1 minute
ConditionVirtualization=no

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
