# examples/destiny/node-exporter/templates/node-exporter-nvme.timer.tmpl
{{/*
  Timer for the NVMe collector. Static (vars: {}), .tmpl is required for core.file.rendered.
  Same ConditionPathExistsGlob=/dev/nvme[0-9]* as in .service: guards against a manual start
  without an NVMe device.
*/}}
[Unit]
Description=Run NVMe metrics collection every 1 minute
ConditionPathExistsGlob=/dev/nvme[0-9]*

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
