# examples/destiny/node-exporter/templates/node-exporter-nvme.timer.tmpl
{{/*
  Таймер NVMe-сборщика. Статичен (vars: {}), .tmpl обязателен для core.file.rendered.
  Та же ConditionPathExistsGlob=/dev/nvme[0-9]*, что и в .service: защита от ручного start
  без NVMe-устройства.
*/}}
[Unit]
Description=Run NVMe metrics collection every 1 minute
ConditionPathExistsGlob=/dev/nvme[0-9]*

[Timer]
OnBootSec=0
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
