# examples/destiny/destiny-node-exporter/templates/node_exporter.service.tmpl
{{/*
  systemd-unit node_exporter. Контекст рендера (docs/templating.md §3.2):
    .vars  — явно переданные core.file.rendered (bin_dir, listen, extra_args);
    .self  — стабильные факты хоста (здесь не используются).
  Движок — Go text/template + sprig-allowlist, strict-mode (ADR-010).

  Сервис-аккаунт — DynamicUser=yes (прод-конвенция §2, гибрид-правило): systemd сам
  заводит запертого transient-юзера на время жизни процесса. node_exporter —
  stateless-демон без owned data-dir, стабильный uid ему не нужен, ручные
  core.user/core.group НЕ заводятся. Под root не работает.

  systemd-hardening (прод-конвенция §3, «безопасность на первом месте»,
  requirements.md): NoNewPrivileges + ProtectSystem=strict (весь / только-чтение —
  node_exporter лишь читает /proc и /sys, ничего не пишет, ReadWritePaths не нужен)
  + изоляция home/tmp/ядра/cgroups + RestrictAddressFamilies до AF_INET/AF_INET6
  (слушает TCP, unix-сокет не использует) + пустой CapabilityBoundingSet (никаких
  capabilities — публичные метрики ядра/ФС привилегий не требуют).

  .vars.extra_args — passthrough-канал для любых флагов node_exporter сверх
  --web.listen-address: коллекторы (--collector.*, --no-collector.*), textfile-dir,
  TLS/basic-auth (--web.config.file=<path>; сам web.yml — отдельный файл, не из
  этой destiny). range дописывает каждый элемент отдельным аргументом ExecStart.
*/}}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
ExecStart={{ .vars.bin_dir }}/node_exporter --web.listen-address={{ .vars.listen }}{{ range .vars.extra_args }} {{ . }}{{ end }}
Restart=on-failure
RestartSec=5

# systemd-hardening (прод-конвенция §3). node_exporter только читает /proc,/sys —
# ProtectSystem=strict без ReadWritePaths достаточно.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
