# deploy/systemd/keeper.service
# systemd-юнит Keeper (ADR-004) — центральный узел Soul Stack.
#
# Установка из пакета soul-stack-keeper (deb/rpm) — пользователь soul-stack,
# /etc/keeper и /var/lib/soul-stack-keeper создаёт postinstall. Ручные шаги
# оператора после установки:
#   install -m0640 keeper.yml /etc/keeper/keeper.yml
#   keeper init --archon=<aid>            # bootstrap первого Архонта (ADR-013)
#   systemctl enable --now keeper
#
# Ручная установка без пакета (для справки):
#   useradd --system --no-create-home --shell /usr/sbin/nologin soul-stack   # один раз на хост
#   install -m0755 keeper/bin/keeper /usr/local/bin/keeper
#   install -d -o soul-stack -g soul-stack /etc/keeper /var/lib/soul-stack-keeper
#   install -m0640 keeper.yml /etc/keeper/keeper.yml
#   install -m0644 deploy/systemd/keeper.service /etc/systemd/system/keeper.service
#   install -m0644 deploy/systemd/keeper.env     /etc/keeper/keeper.env
#   systemctl daemon-reload && systemctl enable --now keeper

[Unit]
Description=Soul Stack Keeper
Documentation=https://github.com/souls-guild/soul-stack
# Keeper слушает gRPC/HTTP и ходит в Postgres/Redis — нужна поднятая сеть.
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=soul-stack
Group=soul-stack

# Путь к конфигу вынесен в EnvironmentFile — оператор меняет его без правки юнита.
EnvironmentFile=/etc/keeper/keeper.env
ExecStart=/usr/local/bin/keeper run --config=${KEEPER_CONFIG}

# Падение демона — перезапуск; не зацикливаемся при битом конфиге (StartLimit).
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=5

# --- hardening (безопасность на первом месте) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
# Единственный writable-путь — рантайм-корень keeper (кеши plugins/services/
# destiny/plugin-src; см. keeper main.go defaults). Должен совпадать с каталогом,
# который создаёт postinstall с владельцем soul-stack, иначе strict-hardening
# заблокирует запись и демон упадёт.
ReadWritePaths=/var/lib/soul-stack-keeper

[Install]
WantedBy=multi-user.target
