# examples/destiny/destiny-redis-exporter/templates/redis_exporter.service.tmpl
{{/*
  systemd-unit redis_exporter. Контекст рендера (docs/templating.md §3.2):
    .vars  — явно переданные core.file.rendered (bin_dir, listen, redis_addr,
             password, extra_args);
    .self  — стабильные факты хоста (здесь не используются).
  Движок — Go text/template + sprig-allowlist, strict-mode (ADR-010).

  Подключение к Redis — через UNIX-СОКЕТ: --redis.addr=unix:///<path>
  (.vars.redis_addr уже содержит префикс unix://, собран в tasks/main.yml из
   input.redis_socket — тот же путь, что unixsocket в redis.conf destiny redis-single).

  Сервис-аккаунт — DynamicUser=yes (прод-конвенция §2, гибрид-правило): systemd сам
  заводит запертого transient-юзера на время жизни процесса. redis_exporter —
  stateless-демон без owned data-dir, стабильный uid ему не нужен, ручные
  core.user/core.group НЕ заводятся. Под root не работает.

  Group-доступ к unix-сокету Redis (unixsocketperm 770, owner redis:redis): сокет
  пускает «группу» (rwx), поэтому transient-юзеру нужно членство в группе redis.
  Оно даётся через SupplementaryGroups=redis — это совместимо с DynamicUser:
  systemd добавляет эфемерный uid в существующую статически-named группу на время
  жизни процесса (группу redis заводит как side-effect пакет redis-server — postinst
  в шаге core.pkg.installed destiny redis-single, применяемой раньше, — а не явный
  core.group-шаг). Прав владельца redis exporter НЕ получает, только членство в группе.

  Пароль Redis передаётся через Environment=REDIS_PASSWORD (redis_exporter читает
  его из env), а не в ExecStart — чтобы не светиться в `ps`/journalctl. Файл unit-а
  пишется mode 0644; на хосте, где это важно, ограничьте доступ к каталогу
  /etc/systemd/system или используйте drop-in с EnvironmentFile mode 0600.

  systemd-hardening (прод-конвенция §3, «безопасность на первом месте»,
  requirements.md): NoNewPrivileges + ProtectSystem=strict (весь / только-чтение)
  + изоляция home/tmp/ядра/cgroups + RestrictAddressFamilies. В отличие от
  node-exporter, тут AF_UNIX ОБЯЗАТЕЛЕН в списке: redis_exporter ходит в Redis
  через unix-сокет, без AF_UNIX connect() к сокету упрётся в EACCES/блокировку
  семейства адресов (прод-конвенция §3 — AF_UNIX добавляется только реально
  использующим unix-сокет демонам). Расширенный блок (MemoryDenyWriteExecute,
  SystemCallFilter=@system-service, RestrictNamespaces, LockPersonality) дополняет
  базовый: exporter — статический Go-бинарь без JIT и без работы с namespaces.

  .vars.extra_args — passthrough-канал для любых флагов redis_exporter сверх
  --redis.addr/--web.listen-address: дополнительные коллекторы (--check-keys=...),
  TLS/basic-auth (--web.config.file=<path>; сам web.yml — отдельный файл, не из
  этой destiny). range дописывает каждый элемент отдельным аргументом ExecStart.
*/}}
[Unit]
Description=Prometheus Redis Exporter
Documentation=https://github.com/oliver006/redis_exporter
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
DynamicUser=yes
# group-доступ к unix-сокету Redis (770) для transient-юзера DynamicUser.
SupplementaryGroups=redis
Environment=REDIS_PASSWORD={{ .vars.password }}
ExecStart={{ .vars.bin_dir }}/redis_exporter --redis.addr={{ .vars.redis_addr }} --web.listen-address={{ .vars.listen }}{{ range .vars.extra_args }} {{ . }}{{ end }}
Restart=on-failure
RestartSec=5

# systemd-hardening (прод-конвенция §3). AF_UNIX добавлен к AF_INET/AF_INET6 —
# redis_exporter подключается к Redis через unix-сокет (--redis.addr=unix://...).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
RestrictNamespaces=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
