# examples/service/mongo/scenario/create/templates/mongod.service.tmpl
{{/*
  systemd unit for mongod for a binary install (essence.install_method=binary). The distro package
  mongodb-org ships its own mongod.service; with an upstream tarball there isn't one -> create renders
  THIS full unit. Go text/template strict (ADR-010). Context: .vars.bin_dir/conf_dir +
  mongo_user/mongo_group.

  Type=simple + foreground (mongod without processManagement.fork under systemd's control -
  it doesn't fork, systemd tracks the PID directly). ExecStart's --config points to
  the rendered mongod.conf (conf_dir/mongod.conf).

  HARDENING (baseline, PILOT): NoNewPrivileges + ProtectSystem=strict + ReadWritePaths for
  the actual mongod write directories (data_dir/log_dir) + home/tmp restrictions. LimitNOFILE/
  LimitNPROC - MongoDB production notes (64000).
  Extended hardening (SystemCallFilter/MemoryDenyWriteExecute) - next slice (mongod is
  more complex than a static Go binary, needs a careful syscall audit).
*/}}
[Unit]
Description=MongoDB Database Server
Documentation=https://docs.mongodb.org/manual
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ .vars.bin_dir }}/mongod --config {{ .vars.conf_dir }}/mongod.conf
User={{ .vars.mongo_user }}
Group={{ .vars.mongo_group }}
Restart=on-failure
RestartSec=5

# Resource limits (MongoDB production notes).
LimitNOFILE=64000
LimitNPROC=64000
LimitMEMLOCK=infinity

# systemd hardening (baseline, PILOT).
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ .vars.data_dir }} {{ .vars.log_dir }}
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
