# tests/e2e-k8s/manifests/soul/soul.service
[Unit]
Description=Soul agent (Soul Stack L3c)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Starts daemon mode. `soul run --config /etc/soul/soul.yml` reads SoulSeed
# from cfg.Paths.Seed (default /var/lib/soul-stack/seed); SoulSeed must
# exist before the unit starts (the harness invokes `soul init` via a
# separate `kubectl exec` before `systemctl start soul.service`).
ExecStart=/usr/local/bin/soul run --config /etc/soul/soul.yml
# on-failure: if the EventStream drops without a reconnect-loop the unit
# restarts (Soul itself has a reconnect-loop; systemd-restart is the
# fallback safety net).
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
