# examples/destiny/redis/templates/redis-sentinel.service.tmpl
{{/*
  systemd unit for redis-sentinel (sentinel mode). Unlike redis-server (whose unit is shipped
  by the distro package), the sentinel daemon is managed by its OWN unit - the package doesn't always ship
  redis-sentinel.service for the right config. A full unit, not a drop-in. Go text/template strict
  (ADR-010). Context: .vars.sentinel_bin/cli_bin (/usr/bin for package installs, /usr/local/bin for
  binary installs), conf_dir/sentinel_port, redis_user/redis_group, data_dir/log_dir/run_dir (RW paths
  for hardening). This way one template serves both install branches.

  Type=notify + supervised systemd (sentinel.conf).
*/}}
[Unit]
Description=Redis Sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart={{ .vars.sentinel_bin }}/redis-sentinel {{ .vars.conf_dir }}/sentinel.conf --supervised systemd
ExecStop={{ .vars.cli_bin }}/redis-cli -p {{ .vars.sentinel_port }} shutdown nosave
Restart=always
User={{ .vars.redis_user }}
Group={{ .vars.redis_group }}
RuntimeDirectory=redis
RuntimeDirectoryMode=0755

# systemd hardening (symmetric to the redis-server hardening drop-in). conf_dir in RW paths
# is REQUIRED: redis-sentinel rewrites sentinel.conf at runtime (persisting topology:
# master/replicas/epoch/known-sentinels), unlike redis-server on create.
ProtectSystem=strict
ReadWritePaths={{ .vars.data_dir }} {{ .vars.log_dir }} {{ .vars.run_dir }} {{ .vars.conf_dir }}
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
