# examples/destiny/redis/templates/redis-sentinel.service.tmpl
{{/*
  systemd-юнит redis-sentinel (режим sentinel). В отличие от redis-server (его юнит несёт
  дистрибутивный пакет), sentinel-демон управляется СВОИМ юнитом — пакет не всегда заводит
  redis-sentinel.service под нужный конфиг. Полный юнит, не drop-in. Go text/template strict
  (ADR-010). Контекст: .vars.sentinel_bin/cli_bin (/usr/bin при package, /usr/local/bin при
  binary), conf_dir/sentinel_port, redis_user/redis_group, data_dir/log_dir/run_dir (RW-пути
  hardening). Так один шаблон обслуживает обе install-ветки.

  Type=notify + supervised systemd (sentinel.conf).
*/}}
[Unit]
Description=Redis Sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart={{ .vars.sentinel_bin }}/redis-sentinel {{ .vars.conf_dir }}/sentinel.conf --supervised systemd
ExecStop={{ .vars.cli_bin }}/redis-cli -p {{ .vars.sentinel_port }} shutdown nosave
Restart=always
User={{ .vars.redis_user }}
Group={{ .vars.redis_group }}
RuntimeDirectory=redis
RuntimeDirectoryMode=0755

# systemd-hardening УНИФИЦИРОВАН с redis-server drop-in (NIM-97; guard —
# TestSentinelUnit_HardeningParity). conf_dir в ReadWritePaths ОБЯЗАТЕЛЕН:
# redis-sentinel переписывает sentinel.conf в рантайме (персист топологии —
# master/replicas/epoch/known-sentinels). MDWE=yes/UMask=007/LimitNOFILE/AF_UNIX —
# обоснование в hardening.conf.tmpl (AF_UNIX критичен: Type=notify sd_notify через
# unix-сокет).
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ .vars.data_dir }} {{ .vars.run_dir }} {{ .vars.log_dir }} {{ .vars.conf_dir }}
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
UMask=007
CapabilityBoundingSet=
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
