# examples/destiny/redis/templates/redis-server.service.tmpl
{{/*
  systemd-юнит redis-server для binary-установки (install.method=binary). Distro-пакет
  несёт свой redis-server.service; при upstream-tarball-е его нет → install.yml рендерит
  ЭТОТ полный юнит. Go text/template strict (ADR-010). Контекст: .vars.bin_dir/conf_dir +
  redis_user/redis_group.

  ★ Имя файла — redis-server.service (НЕ redis.service): так server.yml hardening drop-in
  (redis-server.service.d/hardening.conf) ложится поверх ЭТОГО юнита, и core.service
  рестарт работает В ОБЕИХ install-ветках без ветвления.

  Type=notify + supervised systemd (симметрично дистрибутивному юниту).

  HARDENING тут НЕ дублируется: ProtectSystem/ReadWritePaths/Capability… доезжают
  отдельным drop-in (hardening.conf.tmpl, server.yml) — общий для обеих веток источник
  истины. Здесь — только базовый запуск + User/Group/RuntimeDirectory.
*/}}
[Unit]
Description=Redis In-Memory Data Store
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart={{ .vars.bin_dir }}/redis-server {{ .vars.conf_dir }}/redis.conf --supervised systemd
ExecStop={{ .vars.bin_dir }}/redis-cli shutdown nosave
Restart=always
User={{ .vars.redis_user }}
Group={{ .vars.redis_group }}
RuntimeDirectory=redis
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
