# examples/destiny/redis/templates/redis-server.service.tmpl
{{/*
  systemd-юнит redis-server для binary-установки (install.method=binary). Distro-
  пакет несёт свой redis-server.service; при upstream-tarball-е его нет → install.yml
  рендерит ЭТОТ полный юнит. ★ Имя файла — redis-server.service (НЕ redis.service):
  так server.yml hardening drop-in (/etc/systemd/system/redis-server.service.d/
  hardening.conf) ложится поверх ЭТОГО юнита, и core.service redis-server рестарт
  работает В ОБЕИХ install-ветках без ветвления.

  Движок — Go text/template + sprig-allowlist, strict-mode (ADR-010). .vars:
    .vars.bin_dir   — каталог бинарей binary-ветки (vars.redis_bin_dir, /usr/local/bin —
                      install.yml разложил redis-server/redis-cli туда core.file.present src:);
    .vars.conf_dir  — conf-каталог (redis.conf лежит тут, vars.redis_conf_dir_eff);
    .vars.redis_user/.vars.redis_group — system-аккаунт сервиса (vars.redis_user/_group).
  .self юнит не использует.

  Type=notify + supervised systemd в redis.conf (redis.conf.tmpl: supervised systemd):
  redis-server остаётся foreground-процессом и сам уведомляет systemd о готовности —
  симметрично дистрибутивному юниту, чтобы поведение совпало с package-веткой.

  HARDENING тут НЕ дублируется: ProtectSystem/ReadWritePaths/Capability… доезжают
  отдельным drop-in (hardening.conf.tmpl, server.yml) поверх этого юнита — общий для
  обеих веток источник истины hardening. Здесь — только базовый запуск + User/Group/
  RuntimeDirectory (то, что несёт дистрибутивный unit), drop-in мержит остальное.
*/}}
[Unit]
Description=Redis In-Memory Data Store
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart={{ .vars.bin_dir }}/redis-server {{ .vars.conf_dir }}/redis.conf --supervised systemd
ExecStop={{ .vars.bin_dir }}/redis-cli shutdown nosave
Restart=always
User={{ .vars.redis_user }}
Group={{ .vars.redis_group }}
RuntimeDirectory=redis
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
